# roles/network/files/vpn.service
[Service]
ExecStart=/usr/sbin/openvpn prezi.ovpn
WorkingDirectory=/home/abesto/vpn